# modules/buildbot_slave/templates/buildslave.service.erb
[Unit]
Description=Buildbot Slave
After=multi-user.target

[Service]
Type=simple
ExecStart= /usr/bin/<%= @service_name %> start <%= @slave_dir %>
ExecStop= /usr/bin/<%= @service_name %> stop <%= @slave_dir %>
ExecReload= /usr/bin/<%= @service_name %> reload <%= @slave_dir %>

User=<%= @username %>
Group=<%= @groupname %>
Restart=always

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
