# modules/buildbot_slave/templates/buildslave.service.erb
[Unit]
Description=Buildbot Slave
After=multi-user.target

[Service]
Type=forking
ExecStart= /usr/bin/<%= @service_name %> start .
ExecStop= /usr/bin/<%= @service_name %> stop .
ExecReload= /usr/bin/<%= @service_name %> reload .

User=<%= @username %>
WorkingDirectory=<%= @slave_dir %>
Restart=on-failure

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
Alias=buildslave.service
